# provisioning/supervisor.d/process-worker.conf
[program:process-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:listen
directory = /var/www/html
autostart=true
autorestart=true
user=www-data
numprocs=2
stdout_logfile=/var/log/process-worker-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/process-worker-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
